# deploy/supervisor.conf
[supervisord]
nodaemon=true
[program:nginx]
command=nginx -c /etc/nginx/nginx.conf

[program:exporter]
command=nginx-prometheus-exporter -nginx.scrape-uri http://127.0.0.1:8000/stub_status
